# matahari-broker.service.in
[Unit]
Description=Matahari Broker is an AMQP broker. It receives, stores, routes and forwards messages using the AMQP protcol.
After=syslog.target network.target
Requires=syslog.target network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
NotifyAccess=all
SysVStartPriority=85
User=qpidd
Group=qpidd

EnvironmentFile=-@sysconfdir@/sysconfig/matahari
EnvironmentFile=-@sysconfdir@/sysconfig/matahari-broker

ExecStart=/usr/sbin/matahari-brokerd $QPIDD_OPTIONS

# Restart options include: no, on-success, on-failure, on-abort or always
Restart=on-failure
