# matahari-vios-proxy-guest.service.in
[Unit]
Description=Matahari Vios Proxy Guest
BindTo=dev-virtio\x2dports-org.matahariproject.vios_proxy.0.device
After=dev-virtio\x2dports-org.matahariproject.vios_proxy.0.device

[Service]
ExecStart=-/usr/bin/vios-proxy-guest
Restart=always
RestartSec=0

[Install]
