# sys/matahari.service.in
[Unit]
Description=Matahari @BASE@@BASE_SUB@ daemon using the AMQP protocol.
After=syslog.target network.target matahari-broker.service
Requires=syslog.target network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
NotifyAccess=main
SysVStartPriority=99

EnvironmentFile=-@sysconfdir@/sysconfig/matahari
ExecStart=/usr/sbin/matahari-qmf-@BASE@@BASE_SUB@d $MATAHARI_@DAEMON_TYPE@_ARGS $@BASE@_@DAEMON_TYPE@_ARGS

# Restart options include: no, on-success, on-failure, on-abort or always
Restart=on-failure
