# sys/matahari-broker.upstart.in
# matahari-broker - matahari qpidd broker daemon
#
# Matahari Broker is an AMQP broker. It receives, stores, routes and 
# forwards messages using the AMQP protcol.

description "Matahari broker"

start on filesystem
stop on runlevel [06]

expect fork
respawn

script
    . @sysconfdir@/default/matahari
    . @sysconfdir@/default/matahari-broker
    . @sysconfdir@/default/matahari-vios-proxy
    
    su -c "matahari-brokerd $QPIDD_OPTIONS" qpidd
end script
